# auto_down.cron
*/1 * * * * /usr/src/OneGoAPI/auto_down.sh
